# srv-arm/run_rsynk.service
[Unit]
Description=Синхронизация профиля student с сервера
After=network-online.target
[Service]
ExecStart=/etc/cod-os/run-sync.sh
[Install]
WantedBy=multi-user.target